# systemd/user/ccdbind.service
[Unit]
Description=CCD bind daemon (user)
Wants=game.slice
After=game.slice

[Service]
Type=simple
ExecStart=%h/.local/bin/ccdbind --config %h/.config/ccdbind/config.toml
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target
