# systemd/user/game.slice
[Unit]
Description=Game slice (user)

[Slice]
CPUAccounting=yes

